# com.corec.tracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.corec.tracker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/krish/corec-tracker/src/fetch.js</string>
    </array>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>WorkingDirectory</key>
    <string>/Users/krish/corec-tracker</string>

    <key>StandardOutPath</key>
    <string>/Users/krish/corec-tracker/data/tracker.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/krish/corec-tracker/data/tracker.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
